# misc/systemd/skyringd.service
[Unit]
Description=skyring server

[Service]
Type=forking
KillMode=process
ExecStart=/etc/init.d/skyringd start
ExecStop=/etc/init.d/skyringd stop
